# applications/luci-app-adm/files/etc/uci-defaults/luci-adm
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@adm[-1]
	add ucitrack adm
	set ucitrack.@adm[-1].init=adm
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
